# imp-mail.service
echo "🔧 Fixing Docker container restart issues..."
[Unit]
Description=Imp Mail Application
After=network.target
StartLimitInterval=0

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/ec2-user/imp_mail
ExecStart=/usr/local/bin/docker-compose -f docker-compose.prod.yml -f docker-compose.aws.yml up -d
ExecStop=/usr/local/bin/docker-compose -f docker-compose.prod.yml -f docker-compose.aws.yml down
ExecReload=/usr/local/bin/docker-compose -f docker-compose.prod.yml -f docker-compose.aws.yml restart
TimeoutStartSec=0
Restart=on-failure
RestartSec=10
User=ec2-user
Group=ec2-user

[Install]
WantedBy=multi-user.target